# copilot-api-proxy.service
[Unit]
Description=GitHub Copilot API Proxy (OpenAI & Anthropic Compatible)
Documentation=https://github.com/IT-BAER/copilot-api-proxy
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/copilot-api-proxy
Environment="PORT=4141"
Environment="HOST=0.0.0.0"
Environment="TOKEN_FILE=/opt/copilot-api-proxy/data/github_token.json"
Environment="ACCOUNT_TYPE=individual"
ExecStart=/opt/copilot-api-proxy/venv/bin/python -u app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/copilot-api-proxy/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
